# victron-mppt.service
[Unit]
Description=Victron MPPT MQTT Publisher
After=network.target bluetooth.service
Wants=network-online.target
Requires=bluetooth.service

[Service]
Type=simple
User=alex
Group=alex
WorkingDirectory=/home/alex/victron
Environment=PATH=/home/alex/.local/bin:/home/alex/victron/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/alex/.local/bin/uv run main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal


# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/alex/victron /home/alex/.cache

[Install]
WantedBy=multi-user.target
